# crawler/ots-crawler.service
[Unit]
Description=OpenTrustSeal Crawler Fetch Service
After=network.target

[Service]
Type=simple
User=ott
Group=ott
WorkingDirectory=/opt/ots-crawler
EnvironmentFile=/etc/opentrustseal/crawler.env
ExecStart=/opt/ots-crawler/venv/bin/uvicorn fetch_service:app --host 10.120.0.3 --port 8901 --workers 1
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Chromium is a memory-hungry beast; cap it so a runaway page can't OOM the
# whole box. 768M leaves enough for the service itself and kernel + ufw.
MemoryHigh=640M
MemoryMax=768M

# Hardening: the service only needs its working dir and chromium cache
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/ots-crawler /home/ott/.cache/ms-playwright /tmp
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
